# yestubers.service
[Unit]
Description=YT Cut — YouTube Video Downloader & Cutter SaaS
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/yestubers
Environment="PATH=/opt/yestubers/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="SECRET_KEY=change-me-to-random-64-chars"
Environment="STRIPE_PUBLIC_KEY=pk_live_..."
Environment="STRIPE_SECRET_KEY=sk_live_..."
Environment="STRIPE_WEBHOOK_SECRET=whsec_..."
ExecStart=/opt/yestubers/venv/bin/python3 /opt/yestubers/main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
