# yestubers.service
[Unit]
Description=YesTubers — YouTube Video Downloader & Cutter SaaS
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ytcut
EnvironmentFile=/opt/ytcut/.env
Environment="PATH=/opt/ytcut/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/opt/ytcut/venv/bin/python3 /opt/ytcut/main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
