# ytcut.service
[Unit]
Description=YT Cut — YouTube Video Downloader & Cutter SaaS
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ytcut
Environment="SECRET_KEY=change-me-to-random-64-chars"
Environment="STRIPE_PUBLIC_KEY=pk_live_..."
Environment="STRIPE_SECRET_KEY=sk_live_..."
Environment="STRIPE_WEBHOOK_SECRET=whsec_..."
ExecStart=/opt/ytcut/venv/bin/python3 /opt/ytcut/main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
